# systemd/system/systemd-timesyncd.service.d/wsl.conf
# Enable timesyncd on WSL machines
# so WSL clock is synced on resume from suspend of the host.

[Unit]
ConditionVirtualization=
ConditionVirtualization=|!container
ConditionVirtualization=|wsl
