# systemd/system/multipathd.socket.d/container.conf
# multipathd.service don't start on container. We should have the socket not starting too or it will fail.
# This needs to be upstreamed.

[Unit]
ConditionVirtualization=!container
